# geosub-backend/deploy/linux-arm64/systemd/geosub-event-retention.timer
[Unit]
Description=Apply GeoSub raw analytics retention every day

[Timer]
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec=20min
Persistent=true
Unit=geosub-event-retention.service

[Install]
WantedBy=timers.target
