# geosub-backend/deploy/linux-arm64/systemd/geosub-db-backup.timer
[Unit]
Description=Create a verified GeoSub database backup every day

[Timer]
OnCalendar=*-*-* 02:40:00
RandomizedDelaySec=20min
Persistent=true
Unit=geosub-db-backup.service

[Install]
WantedBy=timers.target
