# geosub-backend/deploy/linux-arm64/systemd/geosub-exchange-rate-sync.timer
[Unit]
Description=Run GeoSub exchange-rate sync every 12 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=12h
RandomizedDelaySec=10min
Persistent=true
Unit=geosub-exchange-rate-sync.service

[Install]
WantedBy=timers.target
