# geosub-backend/deploy/linux-arm64/systemd/geosub-price-pipeline.timer
[Unit]
Description=Run GeoSub price pipeline daily

[Timer]
OnCalendar=*-*-* 03:35:00
RandomizedDelaySec=20m
Persistent=true

[Install]
WantedBy=timers.target
